# tools/prometheus/rustchain-exporter.service
[Unit]
Description=RustChain Prometheus Exporter
After=network.target
Documentation=https://github.com/Scottcjn/Rustchain

[Service]
Type=simple
User=rustchain
Group=rustchain
WorkingDirectory=/opt/rustchain-exporter
ExecStart=/usr/bin/python3 /opt/rustchain-exporter/rustchain_exporter.py
Restart=always
RestartSec=10

# Environment variables
Environment="RUSTCHAIN_NODE_URL=https://rustchain.org"
Environment="EXPORTER_PORT=9100"
Environment="SCRAPE_INTERVAL=60"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/rustchain-exporter

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rustchain-exporter

[Install]
WantedBy=multi-user.target
